# supervisor.conf
# Supervisor 配置文件
# 用于管理工业 AI 调研 Agent 的后台服务

# API 服务配置
[program:industrial-ai-api]
command=/root/工业AI调研agent/venv/bin/python /root/工业AI调研agent/app.py
directory=/root/工业AI调研agent
autostart=true
autorestart=true
stderr_logfile=/var/log/industrial-ai-api.err.log
stdout_logfile=/var/log/industrial-ai-api.out.log
user=root
environment=PATH="/root/工业AI调研agent/venv/bin"
stopsignal=INT
stopwaitsecs=10

# Web 界面配置
[program:industrial-ai-web]
command=/root/工业AI调研agent/venv/bin/python /root/工业AI调研agent/web_ui.py
directory=/root/工业AI调研agent
autostart=true
autorestart=true
stderr_logfile=/var/log/industrial-ai-web.err.log
stdout_logfile=/var/log/industrial-ai-web.out.log
user=root
environment=PATH="/root/工业AI调研agent/venv/bin"
stopsignal=INT
stopwaitsecs=10